# package_scripts/postinstall
#!/bin/bash

/bin/launchctl load /Library/LaunchDaemons/com.github.ygini.UpdateNameFromGLPI.plist

exit 0
